# private_Library/LaunchAgents/com.domus.sort.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.domus.sort</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{ .chezmoi.homeDir }}/.local/bin/domus-sort</string>
        <string>--watch</string>
    </array>

    <!-- Start at login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Keep running -->
    <key>KeepAlive</key>
    <true/>

    <!-- Background process settings -->
    <key>ProcessType</key>
    <string>Background</string>

    <!-- Nice value (lower priority) -->
    <key>Nice</key>
    <integer>15</integer>

    <!-- Low I/O priority -->
    <key>LowPriorityIO</key>
    <true/>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>{{ .chezmoi.homeDir }}/.local/state/domus/launchd-sort-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>{{ .chezmoi.homeDir }}/.local/state/domus/launchd-sort-stderr.log</string>

    <!-- Environment -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:{{ .chezmoi.homeDir }}/.local/bin</string>
        <key>HOME</key>
        <string>{{ .chezmoi.homeDir }}</string>
    </dict>

    <!-- Throttle restarts -->
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
